# crontab
# minute hour day_of_month month day_of_week command
*/5 * * * * source $HOME/.Xdbus; vdirsyncer sync 
*/5 * * * * ~/.bin/battery-warning.sh
